# sources/novaforge_sources/products/release/3.8.0/process/sonar/5.6.3/install/2.1/src/main/resources/bin/sonar.service
[Unit]
Description=Sonar Quality Server
Requires=@MARIADB_SERVICE@.service
After=@MARIADB_SERVICE@.service

[Service]
Type=oneshot

## User configuration
User=@NOVAFORGE_USER@
Group=@NOVAFORGE_GROUP@

## Scripts
ExecStart=@SONAR_BIN@ start 
ExecStop=@SONAR_BIN@ stop
RemainAfterExit=yes

## Give a reasonable amount of time for the server to start up/shut down
TimeoutSec=300

## Place temp files in a secure directory, not /tmp
PrivateTmp=true

[Install]
WantedBy=multi-user.target

